# iso.service
[Unit]
Description=iso image generator
Wants=iso.timer

[Service]
User=root
Type=simple
ExecStart=/usr/share/lilac/iso/gen_iso.sh

